# VTrade.Framework/applications/mql5/Experts/Grande/mcp/analyze_sentiment_server/grande-sentiment.service
[Unit]
Description=Grande Sentiment MCP Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=mcpuser
Group=mcpuser
WorkingDirectory=/opt/grande-sentiment
Environment=MCP_TRANSPORT=streamable-http
Environment=SENTIMENT_PROVIDER=openai_compat
Environment=OPENAI_BASE_URL=http://localhost:11434/v1
Environment=OPENAI_API_KEY=EMPTY
Environment=OPENAI_MODEL=gpt-4o-mini
ExecStart=/opt/grande-sentiment/.venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grande-sentiment

# Health check
ExecStartPost=/bin/sleep 5
ExecStartPost=/bin/bash -c 'for i in {1..30}; do curl -f http://localhost:8000/health && break || sleep 2; done'

[Install]
WantedBy=multi-user.target
